# deploy/tunnelatlas.service
[Unit]
Description=TunnelAtlas reporting daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tunnelatlasd run
Restart=on-failure
RestartSec=5s
User=root
StateDirectory=tunnelatlas
ConfigurationDirectory=tunnelatlas
WorkingDirectory=/var/lib/tunnelatlas
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
